# contrib/init/diminutivecoin.conf
description "Diminutivecoin Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env DIMINUTIVECOIND_BIN="/usr/bin/diminutivecoind"
env DIMINUTIVECOIND_USER="diminutivecoin"
env DIMINUTIVECOIND_GROUP="diminutivecoin"
env DIMINUTIVECOIND_PIDDIR="/var/run/diminutivecoind"
# upstart can't handle variables constructed with other variables
env DIMINUTIVECOIND_PIDFILE="/var/run/diminutivecoind/diminutivecoind.pid"
env DIMINUTIVECOIND_CONFIGFILE="/etc/diminutivecoin/diminutivecoin.conf"
env DIMINUTIVECOIND_DATADIR="/var/lib/diminutivecoind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # bitcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$DIMINUTIVECOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run diminutivecoind."
        echo "The setting must appear in $DIMINUTIVECOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Diminutivecoin Alert\"" \
            "team@diminutivecoin.com"
        echo
        exit 1
    fi

    mkdir -p "$DIMINUTIVECOIND_PIDDIR"
    chmod 0755 "$DIMINUTIVECOIND_PIDDIR"
    chown $DIMINUTIVECOIND_USER:$DIMINUTIVECOIND_GROUP "$DIMINUTIVECOIND_PIDDIR"
    chown $DIMINUTIVECOIND_USER:$DIMINUTIVECOIND_GROUP "$DIMINUTIVECOIND_CONFIGFILE"
    chmod 0660 "$DIMINUTIVECOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$DIMINUTIVECOIND_PIDFILE" \
    --chuid $DIMINUTIVECOIND_USER:$DIMINUTIVECOIND_GROUP \
    --exec "$DIMINUTIVECOIND_BIN" \
    -- \
    -pid="$DIMINUTIVECOIND_PIDFILE" \
    -conf="$DIMINUTIVECOIND_CONFIGFILE" \
    -datadir="$DIMINUTIVECOIND_DATADIR" \
    -disablewallet \
    -daemon

